# roles/worker/templates/gunicorn.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream={{ gunicorn_address }}
BindIPv6Only=both

[Install]

WantedBy=sockets.target
